# advent_calendar/2021/nim/sample_2.conf
[supervisord]
nodaemon=true

[inet_http_server]
port=127.0.0.1:8000

[program:sample_1]
command=/tmp/sample_1
stopwaitsecs=30

[program:sample_2]
command=/tmp/sample_2
stopwaitsecs=30
